# supervisord.conf
[program:node]
command = node tty.js/lib/tty.js > /home/dotcloud/log.txt
directory = /home/dotcloud/current
