# supervisord.conf
[program:node]
command = node server.js > /home/dotcloud/log.txt
directory = /home/dotcloud/current
