# supervisord.conf
[program:node]
command = node tty.js/lib/tty.js
directory = /home/dotcloud/current
